# deploy/systemd/upstreamwx-api.service
# UpstreamWX briefing API — always-on backend (PRD §7, roadmap §M0.1.1).
#
# Template: bootstrap.sh renders the __TOKENS__ from deploy/config.env and installs the
# result to /etc/systemd/system/upstreamwx-api.service. Runs ONE uvicorn process: the
# app's lifespan starts the in-process SREF-cycle refresh scheduler (api/app.py), which
# is a singleton — do NOT add workers, that would multiply the refresh passes (FR-12).
#
# uvicorn binds loopback only; nginx terminates TLS and proxies to it (see deploy/nginx).

[Unit]
Description=UpstreamWX briefing API (FastAPI + SREF-cycle refresh scheduler)
Documentation=https://github.com/osh-labs/upstreamwx/blob/main/deploy/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__GROUP__
WorkingDirectory=__APP_DIR__
EnvironmentFile=__ENV_FILE__
# Playwright's headless Chromium binary location (server-side PDF export, FR-27).
# deploy.sh installs the binary here via `playwright install chromium`.
Environment=PLAYWRIGHT_BROWSERS_PATH=__APP_DIR__/.playwright-browsers
ExecStart=__APP_DIR__/.venv/bin/uvicorn upstreamwx.api.app:app --host __BIND_HOST__ --port __BIND_PORT__
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
# Give a clean cycle: the scheduler swallows per-cycle errors, so SIGTERM is enough.
KillSignal=SIGTERM
TimeoutStopSec=30

# --- Hardening: this process only needs its code (ro) and its data dir (rw) ----------
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=__DATA_DIR__ __APP_DIR__/.playwright-browsers
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
